# ai_service/deploy/photos-face-ai.service
[Unit]
Description=PhotOS Face AI Worker
After=network-online.target redis-server.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/photos/ai_service
EnvironmentFile=/etc/photos-face-ai.env
ExecStart=/var/www/photos/ai_service/.venv/bin/python /var/www/photos/ai_service/main.py
Restart=always
RestartSec=5
RuntimeMaxSec=3600
MemoryHigh=1200M
MemoryMax=1500M
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

