# ai_service/deploy/photos-laravel-queue.service
[Unit]
Description=PhotOS Laravel Queue Worker
After=network-online.target redis-server.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/photos
ExecStart=/usr/bin/php artisan queue:work --queue=photos,default --sleep=1 --tries=3 --timeout=120 --max-jobs=250 --max-time=3600
Restart=always
RestartSec=5
RuntimeMaxSec=3700
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

