# ai_service/deploy/photos-face-ai-results.service
[Unit]
Description=PhotOS Laravel Face AI Results Consumer
After=network-online.target redis-server.service redis.service
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/photos
ExecStart=/usr/bin/php artisan face-ai:consume-results --timeout=5
Restart=always
RestartSec=5
RuntimeMaxSec=3600
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target

