# ai_service/deploy/supervisor-photos-face-ai.conf
[program:photos-face-ai-worker]
command=/var/www/photos/ai_service/.venv/bin/python /var/www/photos/ai_service/main.py
directory=/var/www/photos/ai_service
user=www-data
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/var/log/photos-face-ai/worker.log
stderr_logfile=/var/log/photos-face-ai/worker-error.log
environment=FACE_AI_REDIS_URL="redis://127.0.0.1:6379/0",FACE_AI_TASK_QUEUE="face-ai:tasks",FACE_AI_RESULT_QUEUE="face-ai:results",FACE_AI_BRAND_DETECTOR="disabled"

[program:photos-face-ai-results]
command=/usr/bin/php /var/www/photos/artisan face-ai:consume-results --timeout=5
directory=/var/www/photos
user=www-data
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/var/log/photos-face-ai/results.log
stderr_logfile=/var/log/photos-face-ai/results-error.log

[program:photos-laravel-queue]
command=/usr/bin/php /var/www/photos/artisan queue:work --sleep=1 --tries=3 --timeout=120 --max-jobs=250 --max-time=3600
directory=/var/www/photos
user=www-data
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/var/log/photos-face-ai/queue.log
stderr_logfile=/var/log/photos-face-ai/queue-error.log
